# systemd/pg-backup@hourly.timer
[Unit]
Description=PG backups (%i) [example]
[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
AccuracySec=1min
[Install]
WantedBy=timers.target
